# udev/99-custom-hwmon.rules
ACTION=="add", SUBSYSTEM=="hwmon", ATTRS{name}=="coretemp", SYMLINK+="hwmon_coretemp", RUN+="/bin/sh -c 'ln -s /sys/class/hwmon/%k /dev/hwmon_coretemp'"
ACTION=="add", SUBSYSTEM=="hwmon", ATTRS{name}=="octo", SYMLINK+="hwmon_octo", RUN+="/bin/sh -c 'ln -s /sys/class/hwmon/%k /dev/hwmon_octo'"
ACTION=="add", SUBSYSTEM=="hwmon", ATTRS{name}=="acpitz", SYMLINK+="hwmon_acpitz", RUN+="/bin/sh -c 'ln -s /sys/class/hwmon/%k /dev/hwmon_acpitz'"
ACTION=="add", SUBSYSTEM=="hwmon", ATTRS{name}=="asus", SYMLINK+="hwmon_mobo", RUN+="/bin/sh -c 'ln -s /sys/class/hwmon/%k /dev/hwmon_mobo'"
ACTION=="add", SUBSYSTEM=="hwmon", ATTRS{name}=="iwlwifi_1", SYMLINK+="hwmon_wifi", RUN+="/bin/sh -c 'ln -s /sys/class/hwmon/%k /dev/hwmon_wifi'"
